# supervisor/conf.d/datadog.conf
[program:datadog]
command=datadog-agent run
stdout_logfile=/var/log/datadog/stdout.log
stderr_logfile=/var/log/datadog/stderr.log
