# supervisor/conf.d/emqx.conf
[program:emqx]
command=emqx foreground
stdout_logfile=/var/log/emqx/stdout.log
stderr_logfile=/var/log/emqx/stderr.log
