# supervisor/conf.d/uwsgi.conf
[program:uwsgi]
command=poetry run uwsgi /app/uwsgi/config.ini
stdout_logfile=/var/log/uwsgi/access.log
stderr_logfile=/var/log/uwsgi/error.log
