# supervisor/conf.d/uwsgi.conf
[program:uwsgi]
command=poetry run uwsgi /app/uwsgi/config.ini
